# templates/ext-port.conf
description "Enabling Quantum external networking port"

start on runlevel [2345]

task

script
	ip link set {{ ext_port }} up
	ip link set {{ ext_port }} mtu {{ mtu }}
end script